# transbot.conf
description "Transbot - A slack bot"
author "Bob Brundrett"
  
start on (local-filesystems and net-device-up)
stop on runlevel [06]
  
respawn
exec /home/transbot/Bots/transbot.py
